# roles/golang/templates/etc/profile.d/golang.sh
export GOPATH=$HOME
export PATH=$PATH:/usr/local/go/bin:$HOME/bin
